# etc/systemd/system/restic-check@local.timer
[Unit]
Description=Check restic backup SFTP for errors on a schedule

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=6h

[Install]
WantedBy=timers.target
