# etc/systemd/system/restic-backup@local.timer
[Unit]
Description=Backup with restic on schedule

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
