# etc/systemd/system/restic-backup@offsite.timer
[Unit]
Description=Backup with restic on schedule

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=6h

[Install]
WantedBy=timers.target
